# deploy/trading-assistant-backup.timer
[Unit]
Description=Daily OpenStock decision desk backup

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
